# rootfs/etc/s6-overlay/s6-rc.d/baidunetdisk/run
#!/command/with-contenv bash

export HOME=/app

s6-setuidgid app /opt/baidunetdisk/baidunetdisk --no-sandbox