# rootfs/etc/s6-overlay/s6-rc.d/fluxbox/run
#!/command/with-contenv bash

export HOME=/app
export LANG=zh_CN.UTF-8

exec s6-setuidgid app fluxbox