# rootfs/etc/s6-overlay/s6-rc.d/websockify/run
#!/bin/bash
s6-setuidgid app websockify --web=/opt/noVNC 8080 localhost:5900